# ops/launchd/com.lathe.orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.lathe.orchestrator</string>

	<!-- 5 分間隔で 1 パス実行（orchestrator は常駐せず lock で二重起動を防ぐ）。
	     script は絶対パス必須: 相対だと orchestrator.mjs の isMain 判定
	     （fileURLToPath === argv[1]）が偽になり main が走らない。 -->
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/node</string>
		<string>/Users/cherie/LLMWiki/projects/lathe/scripts/orchestrator.mjs</string>
		<string>--max</string>
		<string>5</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/cherie/LLMWiki/projects/lathe</string>

	<key>StartInterval</key>
	<integer>300</integer>

	<key>RunAtLoad</key>
	<false/>

	<!-- launchd は login shell を継がない — 子プロセスが解決するコマンドの実パスを明示。
	     node / gh / pnpm = /opt/homebrew/bin、claude = ~/.local/bin（2026-07-07 実測）。 -->
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/Users/cherie/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<!-- launchd は同一ファイルへ追記（O_APPEND）。dir 生成と 7 日超 rotate は
	     install.sh と orchestrator（beginPassLog）が担う。 -->
	<key>StandardOutPath</key>
	<string>/Users/cherie/LLMWiki/projects/lathe/.lathe/logs/orchestrator.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/cherie/LLMWiki/projects/lathe/.lathe/logs/orchestrator.log</string>
</dict>
</plist>
