# ops/systemd/lathe-orchestrator.service
[Unit]
Description=Lathe Orchestrator — one-shot pass (#234)
# launchd StartInterval 300 の systemd 等価。timer が 5 分間隔で本 unit を起動する。
# SocketActivation 不使用・podman.socket は orchestrator が必要な場合だけ自己起動。
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

# ── 実行 ──────────────────────────────────────────────────────────────────
# {{NODE_BIN}} = case-setup.sh が nix build/path-info で解決した絶対パス
#   （nix store hash は NixOS rebuild で変化するため直書き禁止、case-setup.sh に委譲）
# 絶対パス必須: orchestrator.mjs は fileURLToPath(import.meta.url) === argv[1] で
#   isMain を判定するため相対パスだと main が走らない（launchd plist コメントと同じ理由）。
# {{REPO_ROOT}} = case-setup.sh がスクリプト位置 ../.. から動的解決（ベタ焼きしない）。
ExecStart={{NODE_BIN}} {{REPO_ROOT}}/scripts/orchestrator.mjs --max 5

WorkingDirectory={{REPO_ROOT}}

# ── 環境変数 ──────────────────────────────────────────────────────────────
# NixOS user systemd は login shell PATH を継がない。
# {{LOGIN_PATH}} = case-setup.sh が install 時に login-shell の $PATH を capture して置換
#   （recon は完全 PATH 文字列を未記録のため guess をコミットせず setup.sh に委譲）。
# node 自体は ExecStart に絶対パスを使うためここに不要。
# gh / claude は ~/.nix-profile/bin 経由でアクセス可（nix profile install 済み前提）。
Environment=PATH={{LOGIN_PATH}}

# ── ログ ─────────────────────────────────────────────────────────────────
# launchd StandardOutPath/StandardErrorPath の等価。
# `append` モード = launchd の O_APPEND と同じ追記挙動。
# ディレクトリと rotate は case-setup.sh と orchestrator (beginPassLog) が担う。
# {{LOG_DIR}} = case-setup.sh が {{REPO_ROOT}}/.lathe/logs として動的解決。
StandardOutput=append:{{LOG_DIR}}/orchestrator.log
StandardError=append:{{LOG_DIR}}/orchestrator.log

# ── 単一インスタンス保証 ───────────────────────────────────────────────────
# Type=oneshot: service が exit するまで次の timer 発火は pending になる（並列起動不可）。
# RefuseManualStart 不採用（plan §2）: 単一インスタンスは Type=oneshot + orchestrator.lock
#   で担保する。`systemctl --user start lathe-orchestrator.service` での即時実行も許容。
